# scripts/aegis-dashboard.service
[Unit]
Description=Aegis Safety Kernel Web Dashboard
Documentation=https://github.com/justinlooney/aegis
After=network.target aegis-verifier.service
Wants=aegis-verifier.service

[Service]
User=aegis
Group=aegis
ExecStart=/usr/bin/python3 -m http.server 8091 --directory /usr/local/share/aegis/dashboard
WorkingDirectory=/usr/local/share/aegis/dashboard
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/usr/local/share/aegis/dashboard
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aegis-dashboard

[Install]
WantedBy=multi-user.target
